# scripts/systemd/hackernel.service
[Unit]
Description=Hackernel

[Service]
Restart=always
RestartSec=3
ExecStart=/usr/sbin/hackernel

[Install]
WantedBy=multi-user.target
