# scripts/systemd/hackernel.service
[Unit]
Description=Hackernel

[Service]
ExecStart=/usr/sbin/hackernel
Nice=-20

[Install]
WantedBy=multi-user.target
