# apps/init/hackernel-notify.service
[Unit]
Description=Hackernel Notify

[Service]
ExecStart=/usr/bin/hackernel-notify
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
